# crontab
* * * * *  SHELL=/usr/bin/bash /usr/local/bin/python /app/update_dns.py >> /var/log/cron.log 2>&1
